# com.wsoftpro.upwork-crm-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wsoftpro.upwork-crm-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mac8868/.gemini/antigravity/scratch/upwork_login_pyautogui/venv/bin/python3</string>
        <string>/Users/mac8868/.gemini/antigravity/scratch/upwork_login_pyautogui/auto_sync.py</string>
        <string>--query</string>
        <string>IT</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mac8868/.gemini/antigravity/scratch/upwork_login_pyautogui</string>

    <key>StartInterval</key>
    <integer>1200</integer>

    <key>StandardOutPath</key>
    <string>/Users/mac8868/.gemini/antigravity/scratch/upwork_login_pyautogui/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mac8868/.gemini/antigravity/scratch/upwork_login_pyautogui/logs/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
